# .docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.out.log

[program:queue-worker]
command=php /var/www/lead-intelligence-system/artisan queue:work --sleep=3 --tries=3 --max-time=3600
user=www-data
autostart=true
autorestart=true
numprocs=1
priority=10
stopwaitsecs=3600
redirect_stderr=true
stdout_logfile=/var/www/lead-intelligence-system/storage/logs/queue-worker.log

[program:scheduler]
; Runs artisan schedule:run every 60 seconds (mirrors a cron * * * * *)
command=/bin/sh -c "while true; do php /var/www/lead-intelligence-system/artisan schedule:run --no-interaction >> /var/www/lead-intelligence-system/storage/logs/scheduler.log 2>&1; sleep 60; done"
user=www-data
autostart=true
autorestart=true
priority=10
redirect_stderr=true
stdout_logfile=/var/www/lead-intelligence-system/storage/logs/scheduler.log
